# fuck-that.service
[Unit]
Description=FUCK
After=network.target nss-lookup.target

[Service]
ExecStart=/usr/bin/python3 /root/socks5/FuckThat.py

[Install]
WantedBy=multi-user.target
